# etc/rc.d/bliptown_helper
#!/bin/ksh

daemon="/usr/local/www/bliptown/script/bliptown_helper"

rc_bg=YES
rc_reload=NO

. /etc/rc.d/rc.subr

pexp="perl ${daemon}"

rc_cmd $1
