# etc/rc.d/block_php
#!/bin/ksh

daemon="/usr/local/www/bliptown/script/block_php"

rc_bg=YES
rc_reload=NO

. /etc/rc.d/rc.subr

pexp="perl:? $daemon"

rc_cmd $1
