# etc/rc.d/block_php
#!/bin/ksh

daemon="/bin/sh"
daemon_flags="/usr/local/www/bliptown/etc/block_php"

rc_bg=YES
rc_reload=NO

. /etc/rc.d/rc.subr

rc_cmd $1
